# iso/config/includes.chroot/etc/xdg/autostart/lucid-titan-console.desktop
[Desktop Entry]
Type=Application
Name=Titan OS Welcome
Comment=First-boot welcome and system readiness check
Exec=bash /opt/titan/bin/titan-welcome
Icon=dialog-information
Terminal=false
Categories=System;
StartupNotify=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=5
Hidden=false
